# GreenParty.IdentitySync/deploy/systemd/greenparty-identity-sync.timer
[Unit]
Description=Run local Green Party identity sync every 15 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Persistent=true
Unit=greenparty-identity-sync.service

[Install]
WantedBy=timers.target
