# MkGreens.IdentitySync/deploy/systemd/mkgreens-identity-sync.timer
[Unit]
Description=Run MK Greens identity sync every 15 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Persistent=true
Unit=mkgreens-identity-sync.service

[Install]
WantedBy=timers.target
