# MkGreens.IdentitySync/deploy/systemd/mkgreens-identity-sync.service
[Unit]
Description=MK Greens identity sync job
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/mkgreens/identity-sync
Environment=DOTNET_ENVIRONMENT=Production
ExecStart=/usr/bin/dotnet /opt/mkgreens/identity-sync/MkGreens.IdentitySync.dll --once
